# autostart/.config/autostart/printUrl.desktop
[Desktop Entry]
Name=PrintOut-Autostart-Script
Type=Application
Exec=bash printUrl.sh
Terminal=true
