# autostart/.config/autostart/startServer.desktop
[Desktop Entry]
Name=Autostart-Script-Server
Comment=Kommentar
Type=Application
Exec=bash startServer.sh
Terminal=true
