# autostart/.config/autostart/startNgrok.desktop
[Desktop Entry]
Name=Ngrok-Autostart
Type=Application
Exec=bash startNgrok.sh
Terminal=true
