# autostart/.config/autostart/startScanner.desktop
[Desktop Entry]
Name=Scanner-Autostart-Script
Type=Application
Exec=bash startScanner.sh
Terminal=true
